# supervisor/bot.conf
[program:bot]
command=python3 /home/user/monitora/bot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log
