# supervisor/checker.conf
[program:checker]
command=python3 /home/user/monitora/checker.py
autostart=true
autorestart=true
stderr_logfile=/var/log/checker.err.log
stdout_logfile=/var/log/checker.out.log
