# supervisor/endpoint.conf
[program:endpoint]
command=python3 /home/user/monitora/endpoint.py
autostart=true
autorestart=true
stderr_logfile=/var/log/endpoint.err.log
stdout_logfile=/var/log/endpoint.out.log
